# ropieee/res/ropieee1-reboot.timer
[Unit]
Description=Reboot

[Timer]
OnCalendar=03:04

[Install]
WantedBy=timers.target

