# ropieee/systemd/ropieee-update.timer
[Unit]
Description=Run RoPieee-update

[Install]
WantedBy=timers.target

