# ropieee/systemd/ropieee1-update.timer
[Unit]
Description=Run RoPieee-update on 4hrs interval

[Timer]
OnBootSec=15min
OnUnitActiveSec=4h

[Install]
WantedBy=timers.target

