# ropieee/systemd/ropieee-update.timer
[Unit]
Description=Run RoPieee-update

[Timer]
OnBootSec=15min
OnUnitActiveSec=4h

[Install]
WantedBy=timers.target

