# ropieee/systemd/ropieee1-update.timer
[Unit]
Description=Run ropieee-update on every boot

[Timer]
OnBootSec=15min

[Install]
WantedBy=timers.target

